# ansible/roles/pg/templates/barman.cron
21 * * * * barman /etc/barman.d/barman.shell.cron
